# falcond-herald.service
[Unit]
Description=Falcond Herald - Power Profile Notification Daemon
Documentation=https://github.com/MurderFromMars/Falcond-Herald
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/falcond-herald
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
